# fedoraConfig/raywork/rc.local
#!/bin/bash
rfkill unblock wifi
sh /home/lihui/startwifi.sh
